# test-network-restart/test-network.service
[Unit]
Description=Runs /root/test-network.sh every hourish

[Service]
Type=oneshot
ExecStart=/root/test-network.sh
