# lets-encrypt/lets-encrypt.service
[Unit]
Description=Let's Encrypt renewal
 
[Service]
Type=oneshot
ExecStart=/opt/lets-encrypt/custom-certbot.sh
